# src/config/serverstick-agent.service
[Unit]
Description=ServerStick Pi Agent
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/serverstick/src/agent
ExecStart=/opt/serverstick/src/agent/.venv/bin/python -m uvicorn main:app --host 0.0.0.0 --port 8080
Restart=always
RestartSec=5
Environment=SERVERSTICK_DIR=/etc/serverstick
Environment=SERVERSTICK_DATA=/var/lib/serverstick/data
EnvironmentFile=-/etc/serverstick/agent.env
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target